# ToUpload/app-supervisor.conf
[program:middleTier]
command=/usr/local/bin/app
autostart=true
autorestart=true